# scripts/systemd/mailcleaner-stopped.target
[Unit]
Description=Mailcleaner - SMTP Antivirus/Antispam Gateway
Requires=sshd.service
Conflicts=mailcleaner.target
AllowIsolate=yes
