# scripts/systemd/spamd@.service
[Unit]
Description=Perl-based spam filter using text analysis
After=syslog.target network-online.target

[Service]
User=mailcleaner
Group=mailcleaner
Type=simple
PIDFile=/run/spamd.pid
EnvironmentFile=-/etc/default/spamd
ExecReload=/bin/kill -HUP $MAINPID
StandardOutput=null
StandardError=null
Restart=always

ProtectSystem=full
PrivateTmp=true
PrivateDevices=true
ProtectHostname=true
ProtectClock=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictNamespaces=true
LockPersonality=true
RestrictRealtime=true
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
