# scripts/systemd/spamd@newsld.service.d/override.conf
[Unit]
Description=Newsletter Daemon daemon for MailCleaner
After=mariadb@slave.service
Requires=mariadb@slave.service
PartOf=mailcleaner.target

[Service]
Type=forking
PIDFile=/var/mailcleaner/run/spamd/%I.pid
User=mailcleaner
Group=mailcleaner
ExecStartPre=
ExecStartPre=+/usr/mailcleaner/bin/dump_%I_config.pl
ExecStart=
ExecStart=/usr/sbin/spamd --socketpath=/var/mailcleaner/spool/spamd/%I.sock --pidfile=/var/mailcleaner/run/spamd/%I.pid --syslog=/var/mailcleaner/log/mailscanner/%I.log --username=mailcleaner --groupname=mailcleaner --siteconfigpath=/usr/mailcleaner/share/%I --nouser-config --daemonize
ProtectSystem=false
PrivateTmp=false
PrivateDevices=false
ProtectHostname=false
ProtectClock=false
ProtectKernelTunables=false
ProtectKernelModules=false
ProtectKernelLogs=false
ProtectControlGroups=false
RestrictNamespaces=false
LockPersonality=false
RestrictRealtime=false
RestrictSUIDSGID=false
Restart=on-failure

[Install]
WantedBy=mailcleaner.target
